# feedbin/supervisor/feedbin.conf
[program:feedbin]
command=/srv/scripts/feedbin.sh
directory=/srv/apps/feedbin/current
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
