# feedbin/supervisor/memcached.conf
[program:memcached]
command=/srv/scripts/memcached.sh
directory=/srv/apps
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
