# feedbin/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
directory=/var/lib/redis
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
